# setup/com.wordoftheday.wake.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
        NOTE: This file is a TEMPLATE only.
        The install.sh script generates the real plist with correct paths
        for the current machine. Do not load this file directly.
    -->
    <key>Label</key>
    <string>com.wordoftheday.wake</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>INSTALL_DIR_PLACEHOLDER/wake_daemon.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/wordoftheday_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/wordoftheday_stderr.log</string>

    <key>ProcessType</key>
    <string>Background</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/bin:/bin:/usr/local/bin:/opt/homebrew/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
